# board/boundarydevices/common/overlay-weston/etc/profile.d/weston.sh
export XDG_RUNTIME_DIR=/run/user/`id -u ${WESTON_USER}`
